# destiny-backend.service
[Unit]
Description=Destiny Decoder FastAPI Backend
After=network.target

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/var/www/destiny-decoder/backend
Environment="PATH=/var/www/destiny-decoder/backend/venv/bin"
ExecStart=/var/www/destiny-decoder/backend/venv/bin/gunicorn main:app \
    --workers 2 \
    --worker-class uvicorn.workers.UvicornWorker \
    --bind 127.0.0.1:8000 \
    --access-logfile /var/log/destiny-backend/access.log \
    --error-logfile /var/log/destiny-backend/error.log \
    --log-level info
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=5
PrivateTmp=true
Restart=on-failure
RestartSec=5s

[Install]
WantedBy=multi-user.target
